# Dataset/Samples/incident_o3-mini.json
"{\n  \"incident_id\": \"17ed34d4-eb45-4937-a4c7-fe3eac228e11\",\n  \"technique_id\": \"T1053.005\",\n  \"technique_desc\": \"Scheduled Task/Job: Scheduled Task\",\n  \"incident_description\": \"On 2023-10-12, suspicious scheduled task activity was detected on both Server-1 and Workstation-12, indicating potential privilege escalation and lateral movement. The incident involved the creation and modification of a scheduled task named 'DailyCleanup' that executed an anomalous PowerShell script under system privileges.\",\n  \"attack_logs\": [\n    {\n      \"timestamp\": \"2023-10-12T15:42:03Z\",\n      \"host\": \"Server-1\",\n      \"action\": \"Task Creation\",\n      \"details\": \"User 'svc_daily' created a new scheduled task 'DailyCleanup' via Task Scheduler to run a PowerShell script daily at 02:00 AM.\"\n    },\n    {\n      \"timestamp\": \"2023-10-12T15:43:12Z\",\n      \"host\": \"Workstation-12\",\n      \"action\": \"Task Modification\",\n      \"details\": \"User 'svc_daily' modified the 'DailyCleanup' task properties to execute under SYSTEM privileges, bypassing standard user restrictions.\"\n    },\n    {\n      \"timestamp\": \"2023-10-12T15:44:25Z\",\n      \"host\": \"Server-1\",\n      \"action\": \"Task Execution\",\n      \"details\": \"The 'DailyCleanup' scheduled task was executed unexpectedly using SYSTEM account privileges, indicating potential malicious script execution.\"\n    }\n  ],\n  \"ground_truth_mitigations\": [\n    \"Immediately remove the 'DailyCleanup' scheduled task from Server-1 using administrative access to prevent further execution of the unauthorized script.\",\n    \"Audit the activity of the user 'svc_daily' on both Server-1 and Workstation-12 and disable its account until its legitimacy can be verified.\",\n    \"Review and enforce stricter permissions for Task Scheduler modifications on all hosts, particularly on Server-1 and Workstation-12, to prevent unauthorized changes.\"\n  ]\n}"