# Dataset/Samples/prev_generator/incident_gpt-4.5-preview.json
```json
{
  "incident_id": "a7e8addd-3ecf-4bf6-897e-8ac3ed182ef6",
  "technique_id": "T1053.005",
  "technique_desc": "Scheduled Task/Job: Scheduled Task",
  "incident_description": "An attacker compromised a workstation through a phishing email containing malicious macros. After gaining initial access, the attacker established persistence by creating a scheduled task that executes a malicious PowerShell script every day at midnight. This scheduled task allowed the attacker to maintain persistent access and perform reconnaissance activities within the network undetected for several days.",
  "attack_logs": [
    {
      "timestamp": "2023-10-10T14:37:22Z",
      "host": "WORKSTATION-12",
      "action": "Email phishing delivery",
      "details": "User opened phishing email attachment 'Invoice_October.xlsm', executing malicious VBA macros."
    },
    {
      "timestamp": "2023-10-10T14:39:05Z",
      "host": "WORKSTATION-12",
      "action": "Scheduled task created",
      "details": "Scheduled task 'WindowsUpdateCheck' created to run daily at 00:00, executing encoded PowerShell command from '%APPDATA%\\Windows\\update.ps1'."
    },
    {
      "timestamp": "2023-10-11T00:00:02Z",
      "host": "WORKSTATION-12",
      "action": "Scheduled task execution",
      "details": "Scheduled task 'WindowsUpdateCheck' executed, initiating outbound network connection to external IP 198.51.100.23."
    },
    {
      "timestamp": "2023-10-11T00:01:14Z",
      "host": "WORKSTATION-12",
      "action": "Reconnaissance activity",
      "details": "Malicious PowerShell script executed network reconnaissance commands (net user, net group, ipconfig) and exfiltrated results to attacker-controlled IP."
    }
  ],
  "ground_truth_mitigations": [
    "Isolate compromised host 'WORKSTATION-12' from the network immediately to prevent further lateral movement and exfiltration.",
    "Disable and remove malicious scheduled task 'WindowsUpdateCheck' from affected host, and delete associated scripts located in '%APPDATA%\\Windows\\update.ps1'.",
    "Conduct a thorough malware scan and forensic analysis on the compromised host to identify and remove additional malicious artifacts.",
    "Implement endpoint detection and response (EDR) tools configured to detect and alert on suspicious scheduled task creations and abnormal PowerShell activity.",
    "Educate users on phishing threats and enforce security policies to disable macros in Office documents received from external sources."
  ]
}
```